# infrastructure/docker/services/frontend/etc/service/php-fpm/run
#!/bin/sh

exec /usr/sbin/php-fpm8.0 -y /etc/php/8.0/fpm/php-fpm.conf -O
